# brain/launchd/com.hermes.voiceloop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.voiceloop</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/chenyouwei/Hermes_Brain/.venv/bin/python</string>
        <string>-u</string>
        <string>-m</string>
        <string>scripts.stackchan_voice_loop</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/chenyouwei/Hermes_Brain</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/chenyouwei/Hermes_Brain/memory/logs/voice_loop.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/chenyouwei/Hermes_Brain/memory/logs/voice_loop.log</string>
</dict>
</plist>
