# brain/launchd/com.hermes.stackchan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.stackchan</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/USERNAME/Hermes_Brain/scripts/start_stackchan_gateway.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/USERNAME/Hermes_Brain</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/USERNAME/Hermes_Brain/memory/logs/stackchan_gateway.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/Hermes_Brain/memory/logs/stackchan_gateway.log</string>
</dict>
</plist>
